# src/files/power-schedule.cron.tmpl
{{.CronBegin}}
{{.Minute}} {{.Hour}} * * * /sbin/shutdown -h now
0 * * * * {{.ScriptPath}}
{{.CronEnd}}
